# package/xg040g-kvm-support/files/etc/hotplug.d/tty/99-ch9329-link
#!/bin/sh

# TTY hotplug environments vary between kernel/usbserial versions. Rescan the
# stable sysfs view instead of relying on DEVICENAME being populated.
case "${ACTION:-}" in
	add)
		/usr/sbin/ch9329-detect >/dev/null 2>&1 || true
		;;
	remove)
		rm -f /dev/ch9329
		(
			sleep 1
			/usr/sbin/ch9329-detect >/dev/null 2>&1 || true
		) &
		;;
esac
