# package/xg040g-cloud-pxe/files/etc/hotplug.d/block/95-xg040g-cloud-pxe
#!/bin/sh

case "${ACTION:-}" in
	add|mount)
		/usr/sbin/xg040g-cloud-pxe install-assets >/dev/null 2>&1 || true
		;;
esac
